# agentorchestrator.service
# ══════════════════════════════════════════════════════════════
# Multi-Agent Orchestrator · systemd 服务模板
# 部署说明: 复制到 /etc/systemd/system/agentorchestrator.service
#   sudo cp agentorchestrator.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now agentorchestrator
#
# 默认启动官方后端栈（API + orchestrator + dispatch + outbox）。
# ══════════════════════════════════════════════════════════════
[Unit]
Description=Multi-Agent Orchestrator Backend Stack
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=agentorchestrator
Group=agentorchestrator
WorkingDirectory=/opt/agentorchestrator
ExecStart=/opt/agentorchestrator/agentorchestrator.sh start
ExecStop=/opt/agentorchestrator/agentorchestrator.sh stop
ExecReload=/opt/agentorchestrator/agentorchestrator.sh restart

Environment=AGENTORCHESTRATOR_API_HOST=0.0.0.0
Environment=AGENTORCHESTRATOR_API_PORT=38000

Restart=on-failure
RestartSec=5
TimeoutStartSec=120
TimeoutStopSec=60

NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/agentorchestrator/data /opt/agentorchestrator/logs /opt/agentorchestrator/.pids

[Install]
WantedBy=multi-user.target
